# setup/dirapa.target
[Unit]
Description=Dirapa target
StopWhenUnneeded=yes